# companion/usage-publisher.service
[Unit]
Description=Claude usage publisher (feeds rpi-sentinel-display)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust the user and paths to your Raspberry Pi / this repo's checkout.
User=pi
WorkingDirectory=/home/pi/rpi-sentinel/companion
# Secrets (MQTT_HOST, MQTT_USER, MQTT_PASS) live outside the repo so they are
# never committed. Create this file with 600 perms, e.g.:
#   MQTT_HOST=your-broker.example.com
#   MQTT_USER=youruser
#   MQTT_PASS=yourpass
EnvironmentFile=/home/pi/usage-publisher.env
Environment=MQTT_PORT=8883
Environment=MQTT_PREFIX=rpi
Environment=MQTT_TLS=true
Environment=POLL_SECONDS=600
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# Optional caps so the display shows a percentage instead of raw tokens:
# Environment=FIVE_HOUR_LIMIT=800000
# Environment=WEEKLY_LIMIT=200000000
# Token is read from ~/.claude/.credentials.json of the User above. If this
# machine never runs `claude`, use a long-lived token from `claude setup-token`:
# Environment=CLAUDE_OAUTH_TOKEN=...
# Environment=CLAUDE_CREDENTIALS=/home/pi/.claude/.credentials.json
ExecStart=/usr/bin/python3 usage_publisher.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
